# server/listener-ap.service
[Unit]
Description=Listener Wi-Fi Access Point
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/fpp/ChatGPT-fpp-listen/server/listener-ap.sh
ExecStop=/bin/bash -c 'pkill -f hostapd-listener.conf || true; pkill -f listener-dnsmasq.conf || true'

[Install]
WantedBy=multi-user.target
