# server/listener-ap.service
[Unit]
Description=Listener Wi-Fi Access Point
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/fpp/listen-sync/listener-ap.sh
ExecStop=/bin/bash -c 'pkill -f hostapd-listener.conf || true; pkill -f hostapd-show.conf || true; pkill -f listener-dnsmasq.conf || true; pkill -f show-dnsmasq.conf || true; nft delete table inet show_ap 2>/dev/null || true; rm -f /home/fpp/listen-sync/show-rewrite.conf /opt/fpp/www/.htaccess 2>/dev/null || true; systemctl reload apache2 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
